# files/jenkins.cron
0 5,18 * * * bash -x /var/lib/jenkins/jenkins-home/sync_git.sh >> /var/lib/jenkins/sync_git.log 2>&1
